# flask-sqlite-s3-project/ansible/roles/flask-app/templates/flask-app.service.j2
[Unit]
Description=Flask Application with SQLite
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/opt/flask-app
ExecStart=/usr/bin/python3 /opt/flask-app/app.py
Restart=always

[Install]
WantedBy=multi-user.target